# systemd/calendar-digest.timer
[Unit]
Description=Run Calendar Digest every Sunday at 18:00

[Timer]
OnCalendar=Sun *-*-* 18:00:00
Persistent=true

[Install]
WantedBy=timers.target
